# infra/systemd/retail-kiosk-agent.service
# Template: install as /etc/systemd/system/retail-kiosk-agent.service
# Edit User/Group/WorkingDirectory if the target kiosk uses different paths.
[Unit]
Description=Retail Kiosk Agent / Heartbeat and Remote Operations
After=network-online.target retail-kiosk-local-backend.service
Wants=network-online.target retail-kiosk-local-backend.service tailscaled.service
StartLimitIntervalSec=180
StartLimitBurst=5
RequiresMountsFor=/opt/retail-kiosk /var/lib/retail-kiosk /var/log/retail-kiosk

[Service]
Type=simple
EnvironmentFile=/etc/retail-kiosk/kiosk.env
WorkingDirectory=/opt/retail-kiosk/current
User=kiosk
Group=kiosk-runtime
ExecStart=/usr/bin/env node services/kiosk-agent/dist/index.js
Restart=on-failure
RestartSec=10s
TimeoutStartSec=30
TimeoutStopSec=20
KillSignal=SIGTERM
SyslogIdentifier=retail-kiosk-agent
StateDirectory=retail-kiosk
LogsDirectory=retail-kiosk
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/retail-kiosk /var/log/retail-kiosk /etc/retail-kiosk /opt/retail-kiosk
UMask=0027

[Install]
WantedBy=retail-kiosk.target
